# frontend/conf/init/service.conf
description "mysite frontend - express"

start on (filesystem)
stop on runlevel [016]

respawn
setuid www-data
setgid www-data
chdir /var/www/TARGET_NAME/frontend

exec /usr/bin/yarn start -- --port=PORT
